# docker-gitlab/docker-radarr-pk/build/root/etc/services.d/git/run
#!/usr/bin/with-contenv bash

if [[ $(date +%M) == '00' ]]; then 

cd /dev/shm/docker
git add . || echo "could no git add $(date)"
git commit -m "$(date)" || echo "git commit failed"
git push
sleep 60
fi
