# docker-torproxy-customexitnodes/supervisord/tor-countries.conf
[program:tor-country]
command=/usr/bin/setcountry.sh
autostart = true
redirect_stderr = true
stdout_logfile_maxbytes = 10MB
stdout_logfile_backups = 10
